# etc/supervisor/conf.d/02-httpd.conf
[program:httpd]
command=/usr/local/bin/httpd-foreground